# SD_ROOT/wz_mini/etc/rc.local
#!/bin/sh

#This file executes last, after rc.local.d/


